# scripts/elliott-4h.service
[Unit]
Description=Elliott Wave 4H Timeframe Analyzer
After=network.target

[Service]
Type=oneshot
User=rezaserver
WorkingDirectory=/home/rezaserver/ElliottWaveAnalyzer
ExecStart=/bin/bash /home/rezaserver/ElliottWaveAnalyzer/scripts/run_automated_analysis.sh --timeframe 4H
